# scripts/cti-system.service
[Unit]
Description=CTI System - Daemon + Dashboard Automático
After=network.target

[Service]
Type=simple
User=farias
WorkingDirectory=/home/farias/Desktop/CTI
Environment="PATH=/home/farias/Desktop/CTI/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/farias/Desktop/CTI/venv/bin/python /home/farias/Desktop/CTI/process_manager.py --port 8502
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cti-system

# Limites
MemoryLimit=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
